# init/stream-archiver.in
#!/bin/bash
### BEGIN INIT INFO
# chkconfig:         2345 90 60
# Provides:          stream-archiver.js
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: forever running stream-archiver.js
# Description:       stream-archiver.js
### END INIT INFO
#
# initd a node app
# Based on a script posted by https://gist.github.com/jinze at https://gist.github.com/3748766
#

if [ -e /lib/lsb/init-functions ]; then
	# LSB source function library.
	. /lib/lsb/init-functions
fi;

pidFile="/var/run/stream-archiver.pid"
logFile="/var/run/stream-archiver.log"

command="node"
nodeApp="stream-archiver.js"
appDir="@@APP_DIR@@"
foreverApp="forever"
foreverOpts="--uid \"$nodeApp\" --pidFile $pidFile -l $logFile --workingDir $appDir --sourceDir $appDir -d"

start() {
	echo "Starting $nodeApp"

	# Notice that we change the PATH because on reboot
   # the PATH does not include the path to node.
   # Launching forever with a full path
   # does not work unless we set the PATH.
   PATH=/usr/local/bin:$PATH
	export NODE_ENV=production
   #PORT=80
   $foreverApp start $foreverOpts -a -c "$command" $nodeApp
   RETVAL=$?
}

restart() {
	echo "Restarting $nodeApp"
	$foreverApp restart $foreverOpts $nodeApp
	RETVAL=$?
}

stop() {
	echo "Shutting down $nodeApp"
   $foreverApp stop $foreverOpts \"$nodeApp\"
   RETVAL=$?
}

status() {
   echo "Status $nodeApp"
   $foreverApp list $foreverOpts
   RETVAL=$?
}

case "$1" in
   start)
        start
        ;;
    stop)
        stop
        ;;
   status)
        status
       ;;
   restart)
   	restart
        ;;
	*)
       echo "Usage:  {start|stop|status|restart}"
       exit 1
        ;;
esac
exit $RETVAL
